# root/etc/services.d/puma/run
#!/usr/bin/with-contenv bash
cd /opt/docker-readme-sync || exit
exec puma -b unix:///var/run/puma.sock config.ru

