# root/etc/services.d/puma/run
#!/usr/bin/with-contenv bash
cd /opt/github-dockerhub-sync
exec puma config.ru

