# root/etc/services.d/puma/run
#!/usr/bin/with-contenv bash
cd /opt/github-dockerhub-sync
alias phantomjs='/usr/bin/xvfb-run /usr/bin/phantomjs'
exec puma config.ru

