# root/etc/services.d/puma/run
#!/usr/bin/with-contenv bash
cd /opt/docker-readme-sync
exec puma config.ru

